# systemd-unit/surveillance-cleaner.service
[Unit]
Description=Cleans old surveillance videos

[Service]
ExecStart=%h/surveillance-cleaner/target/release/surveillance-cleaner %h/ipcam
